# fix-workspace.cron
#!/bin/sh

find $HOME/workspace -type d -exec chmod 755 "{}" \;
